# s6-overlay/s6-rc.d/cloudcli/run
#!/command/with-contenv sh
cd /workspace
export HOME=/home/claude
export WORKSPACES_ROOT=/workspace
export NODE_OPTIONS="${NODE_OPTIONS:+$NODE_OPTIONS }--no-deprecation"
exec s6-setuidgid claude cloudcli --port 3001
